# dot_dotfiles/services/mount-hdds/mount-disks/mount-disks.service
[Unit]
Description=Montar discos
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mount-disks.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
