# dot_dotfiles/services/automatic-shutdown/apagado.service
[Unit]
Description="Apagado script"

[Service]
ExecStart=/usr/local/bin/apagado.sh
